# overlay/system/usr/lib/systemd/user/mediascanner-2.0.service.d/clear-gst-registry.conf
[Service]
ExecStartPre=/bin/sh -c "grep -q mtkhevcdecoder ~/.cache/gstreamer-1.0/registry.aarch64.bin || rm -f ~/.cache/gstreamer-1.0/registry.aarch64.bin"
